# examples/keep_alive.service
[Unit]
Description=Wake on LAN Monitoring Service

[Service]
ExecStart=/usr/bin/python3 /path/to/your/script.py
Restart=always

[Install]
WantedBy=multi-user.target

