# examples/keep_alive.service
[Unit]
Description=Keep-Alive Monitoring Service

[Service]
ExecStart=/path/to/keep_alive_service.sh
Restart=always

[Install]
WantedBy=multi-user.target

